# build_files/etc/systemd/system/system-build.timer
[Unit]
Description=Daily automatic builds for custom bootc image

[Timer]
OnCalendar=daily
OnBootSec=15min

[Install]
WantedBy=timers.target